# infra/systemd/qujing-restore-check.timer
[Unit]
Description=Verify a Qujing database restore every week

[Timer]
OnCalendar=Sun *-*-* 04:30:00 UTC
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
